# cat_location_daemon.service
[Unit]
Description=Publish cat location updates from tracker events
After=network-online.target rabbitmq-server.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/libexec/tcpproxy/cat_location_daemon \
  --input-uri=amqp://guest:guest@127.0.0.1:5672/%2F \
  --input-queue=cat.location.source \
  --input-exchange=tracker.events.parsed \
  --input-routing-key=tracker.parsed \
  --output-uri=amqp://guest:guest@127.0.0.1:5672/%2F \
  --output-exchange=cat.location \
  --output-routing-key=cat.position \
  --kml=/usr/local/etc/Locations.kml \
  --pidfile=/run/cat_location_daemon/cat_location_daemon.pid
PIDFile=/run/cat_location_daemon/cat_location_daemon.pid
RuntimeDirectory=cat_location_daemon
RuntimeDirectoryMode=0755
Restart=on-failure

[Install]
WantedBy=multi-user.target
